# supervisor/laravel-scheduler.conf
[program:laravel-scheduler]
process_name=%(program_name)s_%(process_num)02d
command=php /home/sibabalwentoyi/payments-app/artisan schedule:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=sibabalwentoyi
numprocs=1
redirect_stderr=true
stdout_logfile=/home/sibabalwentoyi/payments-app/storage/logs/scheduler.log
stopwaitsecs=3600
